# deploy/com.giannisan.agora-scribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.giannisan.agora-scribe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/giannisan/workspace/agora/target/release/scribe</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AGORA_HUB</key>
        <string>http://100.84.87.107:8787</string>
        <key>AGORA_ROOM</key>
        <string>dev</string>
        <key>AGORA_DIRS</key>
        <string>/Users/giannisan/workspace</string>
        <key>HOME</key>
        <string>/Users/giannisan</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/agora-scribe.err</string>
</dict>
</plist>
